# dockerfiles/runit/.log/run
#!/bin/sh

SERVICE=$(basename "$(dirname "$PWD")")
exec awk -v svc="$SERVICE" '{ print "[" svc "] " $0; fflush(); }'
