# dockerfiles/runit/pve-cluster/run
#!/bin/bash

if [[ -n "$TZ" ]]; then
  echo "$TZ" > /etc/timezone
fi

# /run/pve-cluster.pid
exec /usr/bin/pmxcfs -f
